# deploy/systemd/tiktok-rkhunter-scan.service
[Unit]
Description=Weekly rkhunter malware/rootkit scan
OnFailure=tiktok-notify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rkhunter --check --sk --rwo
StandardOutput=journal
StandardError=journal

# Deliberately no hardening directives here, same rationale as
# tiktok-lynis-audit.service (wiki/decisions/2026-07-03-lynis-periodic-audit.md):
# a real rkhunter scan needs broad read access across the whole system
# (binaries, kernel modules, config) -- ProtectSystem/ProtectHome would hide
# exactly what it exists to inspect. Runs as root (no User=).
#
# Unlike lynis, OnFailure=tiktok-notify@ IS wired here: rkhunter --check
# returns non-zero specifically when it finds warnings (verified empirically,
# see wiki/decisions/2026-07-03-rkhunter-malware-scanner.md), so its exit
# code is a real "something's wrong" signal, not just an audit surfacing
# findings for a human to browse.
#
# File-property/hash integrity tests are disabled in /etc/rkhunter.conf.local
# (DISABLE_TESTS += hashes properties attributes) -- that's task #14's
# (AIDE) job, not rkhunter's, to avoid double-coverage.
MemoryMax=512M
CPUQuota=100%
TasksMax=64
