# deploy/systemd/tiktok-slot-render@0000.timer
[Unit]
Description=Render trigger for the 00:00 slot (fires 23:30 Europe/Madrid)

[Timer]
OnCalendar=*-*-* 23:30:00 Europe/Madrid
Persistent=false
Unit=tiktok-slot-render@0000.service

[Install]
WantedBy=timers.target
