# deploy/systemd/tiktok-xvfb.service
[Unit]
Description=Xvfb virtual display :99 for Playwright headed uploads
OnFailure=tiktok-notify@%n.service
After=network.target

[Service]
Type=simple
User=christian
Group=christian
ExecStart=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -nolisten tcp
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Hardening (P0.5, research run 7). Can be the strictest unit in the
# fleet: no app-dir access at all, `-nolisten tcp` means no network
# needed either -- only a unix socket for the X11 render/upload units
# to connect to. Deliberately NO PrivateTmp (same reason as
# render/upload/upload-legacy): the /tmp/.X11-unix/X99 socket this
# service creates must stay visible to those units' real /tmp.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
