# deploy/systemd/tiktok-xvfb.service
[Unit]
Description=Xvfb virtual display :99 for Playwright headed uploads
OnFailure=tiktok-notify@%n.service
After=network.target

[Service]
Type=simple
User=christian
Group=christian
ExecStart=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -nolisten tcp
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Hardening (P0.5, research run 7). `-nolisten tcp` means no network
# needed -- only a unix socket for the X11 render/upload units to
# connect to. Deliberately NO PrivateTmp (same reason as
# render/upload/upload-legacy): the /tmp/.X11-unix/X99 socket this
# service creates must stay visible to those units' real /tmp.
#
# MUST be ProtectSystem=full, not strict -- verified empirically.
# strict makes Xvfb crash ~12s after a clean start with "Could not
# create lock file in /tmp/.tX99-lock" (bisected: ProtectSystem=strict
# alone reproduces it; every other directive is fine in isolation).
#
# A second, independent gotcha found during rollout: even with `full`,
# the real tiktok-xvfb.service unit kept crash-looping on production's
# :99 specifically, while a brand-new unit name with the identical
# config was rock solid on the same display. Root cause: stale
# systemd-internal state (private mount-namespace/cgroup remnants)
# left over from an earlier rapid crash-loop under the broken `strict`
# config -- NOT the display number or this config. `systemctl
# daemon-reexec` cleared it; a normal `daemon-reload` does not. If this
# unit is ever hardened again after another crash-loop, try
# daemon-reexec before assuming the new config is wrong.
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
