# deploy/systemd/tiktok-webapp.service
[Unit]
Description=TikTok control-plane FastAPI + SPA
OnFailure=tiktok-notify@%n.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
EnvironmentFile=/etc/tiktok/environment
Environment=TZ=Europe/Madrid
Environment=PYTHONUNBUFFERED=1
Environment=PLAYWRIGHT_BROWSERS_PATH=/srv/tiktok/app/.playwright
ExecStart=/srv/tiktok/app/.venv/bin/python -u -m uvicorn webapp.backend.app:app --host 0.0.0.0 --port 8765 --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening (P0.5, research run 7). ProtectSystem=strict locks the
# whole filesystem read-only except the paths below: the repo root
# (config.toml's atomic-write-via-tempfile pattern creates the temp
# file in its parent dir, so the dir must be writable, not just the
# file) and /etc/tiktok (the dashboard's env-editor UI writes there).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
ReadWritePaths=/srv/tiktok/app /etc/tiktok

# Resource limits (R3.1). Sized against observed steady-state (~66MB,
# 10 tasks) with real headroom -- a runaway leak or fork-bomb should
# hit these before starving render/upload units on the same host.
MemoryMax=512M
CPUQuota=100%
TasksMax=64
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
