# deploy/systemd/tiktok-slot-render@.service
[Unit]
Description=Render pipeline for slot %i (fires 30 min before publish)
OnFailure=tiktok-notify@%n.service
Requires=tiktok-xvfb.service
After=tiktok-xvfb.service

[Service]
Type=oneshot
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
EnvironmentFile=/etc/tiktok/environment
Environment=DISPLAY=:99
Environment=TZ=Europe/Madrid
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=4
Environment=PLAYWRIGHT_BROWSERS_PATH=/srv/tiktok/app/.playwright
TimeoutStartSec=30min
ExecStart=/srv/tiktok/app/.venv/bin/python -u scripts/slot.py render --instance %i
StandardOutput=journal
StandardError=journal

# Hardening (P0.5, research run 7). ProtectSystem=full (not strict) --
# Chromium/ffmpeg need filesystem flexibility strict would block, and
# the app dir stays normally writable under full (only /usr,/boot,/etc
# go read-only), so no explicit ReadWritePaths needed here. Deliberately
# NO PrivateTmp -- Playwright connects to Xvfb over the /tmp/.X11-unix/X99
# unix socket; an isolated private /tmp would hide that socket and break
# every headed-browser render/upload. NOT adding MemoryDenyWriteExecute
# or SystemCallFilter -- both need real testing against Chromium/ffmpeg
# before enabling (research run 7), not assumed safe here.
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

# Resource limits (R3.1). Generous backstop, not a tuned ceiling -- only
# one real render's peak (1.4GB) is known against a 16GB host; MemoryMax
# leaves 10GB+ headroom for OS/other units rather than anchoring tight
# to that single sample. CPUQuota=600% (6/8 cores) so legitimate
# whisper/ffmpeg multithreading during a real 44-minute render isn't
# throttled. LimitNOFILE deliberately left unset -- Chromium/Playwright
# want systemd's generous default, not a lower explicit cap.
MemoryMax=6G
CPUQuota=600%
TasksMax=512
