# deploy/systemd/tiktok-slot-render@1800.timer
[Unit]
Description=Render trigger for the 18:00 slot (fires 17:30 Europe/Madrid)
Requires=tiktok-slot-render@1800.service

[Timer]
OnCalendar=*-*-* 17:30:00 Europe/Madrid
Persistent=false
Unit=tiktok-slot-render@1800.service

[Install]
WantedBy=timers.target
