# deploy/systemd/tiktok-retention.service
[Unit]
Description=Delete artifacts for posts that reached terminal state

[Service]
Type=oneshot
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
EnvironmentFile=/etc/tiktok/environment
Environment=TZ=Europe/Madrid
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/tiktok/app/.venv/bin/python -u -m pipeline.retention
StandardOutput=journal
StandardError=journal

# Hardening (P0.5, research run 7). Pure local file/DB cleanup, no
# network calls at all — restricted to AF_UNIX only.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=
ReadWritePaths=/srv/tiktok/app/data

# Resource limits (R3.1). Short-lived, lightweight file/DB cleanup.
MemoryMax=256M
CPUQuota=50%
TasksMax=32
LimitNOFILE=2048
