# deploy/systemd/tiktok-slot-upload@1800.timer
[Unit]
Description=Upload trigger for the 18:00 slot

[Timer]
OnCalendar=*-*-* 18:00:00 Europe/Madrid
Persistent=false
Unit=tiktok-slot-upload@1800.service

[Install]
WantedBy=timers.target
