# deploy/systemd/tiktok-secrets-backup.timer
[Unit]
Description=Nightly trigger for the secrets/tokens/cookies backup

[Timer]
OnCalendar=*-*-* 03:15:00 Europe/Madrid
Persistent=true
Unit=tiktok-secrets-backup.service

[Install]
WantedBy=timers.target
