# deploy/systemd/tiktok-secrets-backup.service
[Unit]
Description=Nightly age+rclone backup of secrets/tokens/cookies to R2
OnFailure=tiktok-notify@%n.service

[Service]
Type=oneshot
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
Environment=TZ=Europe/Madrid
ExecStart=/srv/tiktok/app/deploy/scripts/tiktok-secrets-backup.sh
StandardOutput=journal
StandardError=journal

# Hardening (task #7, same profile family as tiktok-notify@ -- read-only
# local access + outbound HTTPS to R2, no local writes at all (pure
# tar|age|rclone pipe, no temp files -- PrivateTmp is a no-op safety net,
# not load-bearing). No ReadWritePaths needed: only reads under /etc/tiktok
# and /srv/tiktok/app/data, both allowed read-only under ProtectSystem=strict.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

# Resource limits (R3.1 style). Trivial short-lived pipe, ~few KB payload.
MemoryMax=128M
CPUQuota=25%
TasksMax=16
LimitNOFILE=2048
