# deploy/systemd/tiktok-confirm.timer
[Unit]
Description=Fire tiktok-confirm every 30 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=tiktok-confirm.service

[Install]
WantedBy=timers.target
