# deploy/systemd/tiktok-notify@.service
[Unit]
Description=Send Telegram alert for %i

[Service]
Type=oneshot
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
EnvironmentFile=/etc/tiktok/environment
Environment=TZ=Europe/Madrid
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/tiktok/app/.venv/bin/python /srv/tiktok/app/scripts/notify.py %i
StandardOutput=journal
StandardError=journal
