# deploy/systemd/tiktok-slot-upload@0000.timer
[Unit]
Description=Upload trigger for the 00:00 slot
Requires=tiktok-slot-upload@0000.service

[Timer]
OnCalendar=*-*-* 00:00:00 Europe/Madrid
Persistent=true
Unit=tiktok-slot-upload@0000.service

[Install]
WantedBy=timers.target
