# deploy/systemd/tiktok-rkhunter-scan.timer
[Unit]
Description=Weekly trigger for the rkhunter malware/rootkit scan

[Timer]
OnCalendar=Sun *-*-* 05:00:00 Europe/Madrid
Persistent=true
Unit=tiktok-rkhunter-scan.service

[Install]
WantedBy=timers.target
