# deploy/systemd/tiktok-confirm.service
[Unit]
Description=TikTok confirm-live scrape (moves posted_under_review -> posted)
OnFailure=tiktok-notify@%n.service

[Service]
Type=oneshot
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
EnvironmentFile=/etc/tiktok/environment
Environment=TZ=Europe/Madrid
Environment=PYTHONUNBUFFERED=1
Environment=PLAYWRIGHT_BROWSERS_PATH=/srv/tiktok/app/.playwright
ExecStart=/srv/tiktok/app/.venv/bin/python -u pipeline/confirm_live.py
StandardOutput=journal
StandardError=journal

# Hardening (P0.5, research run 7). DB updates + logs/ (see
# tiktok-bot.service for why logs/ is separate from data/).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
ReadWritePaths=/srv/tiktok/app/data /srv/tiktok/app/logs
