# deploy/systemd/tiktok-secrets-backup.service
[Unit]
Description=Nightly age+curl backup of secrets/tokens/cookies to R2
OnFailure=tiktok-notify@%n.service
# Persistent=true on the timer replays a missed run at boot; without ordering
# on network-online the catch-up fires before DNS is up and the R2 upload
# fails with "Could not resolve host" (seen after the 2026-07-05 outage boot).
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=christian
Group=christian
WorkingDirectory=/srv/tiktok/app
Environment=TZ=Europe/Madrid
ExecStart=/srv/tiktok/app/deploy/scripts/tiktok-secrets-backup.sh
StandardOutput=journal
StandardError=journal

# Hardening (task #7, same profile family as tiktok-notify@ -- read-only
# local access + outbound HTTPS to R2. PrivateTmp=yes IS load-bearing here:
# the script stages the age-encrypted (never plaintext) archive to a temp
# file under its private /tmp before curl-uploading it (a raw S3 PUT needs
# a known Content-Length, unlike a streamed pipe), auto-torn-down on exit.
# No ReadWritePaths needed beyond that: only reads under /etc/tiktok and
# /srv/tiktok/app/data, both allowed read-only under ProtectSystem=strict.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

# Resource limits (R3.1 style). Trivial short-lived pipe, ~few KB payload.
MemoryMax=128M
CPUQuota=25%
TasksMax=16
LimitNOFILE=2048
