# deploy/systemd/tiktok-slot-render@1200.timer
[Unit]
Description=Render trigger for the 12:00 slot (fires 11:30 Europe/Madrid)

[Timer]
OnCalendar=*-*-* 11:30:00 Europe/Madrid
Persistent=false
Unit=tiktok-slot-render@1200.service

[Install]
WantedBy=timers.target
