# launchd/com.sebastian.tiktok-confirm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sebastian.tiktok-confirm</string>

  <key>WorkingDirectory</key>
  <string>/Users/sebastian/Automated-TikTok-Upload</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/sebastian/Automated-TikTok-Upload/venv/bin/python</string>
    <string>-u</string>
    <string>pipeline/confirm_live.py</string>
  </array>

  <!-- Every 30 minutes, all day. Confirm-live doesn't share the
       upload worker's posting window: posts may sit in `posted_under_review`
       any time of day, and moderation can complete outside 19:00-01:59. -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/sebastian</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/sebastian/Automated-TikTok-Upload/data/logs/confirm_live.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sebastian/Automated-TikTok-Upload/data/logs/confirm_live.stderr.log</string>
</dict>
</plist>
